# devflowstate.service
[Unit]
Description=DevFlowState Application
Documentation=https://github.com/devflowstate/devflowstate
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/devflowstate
ExecStart=/usr/bin/node /var/www/devflowstate/server.js
Restart=always
RestartSec=10
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOST=127.0.0.1

# Security
NoNewPrivileges=true
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=devflowstate

[Install]
WantedBy=multi-user.target
